# conf/supervisor/conf.d/button.conf
[program:watch_button]
command=/home/pi/JoustMania/watch_button.sh
user=root
directory=/home/pi/JoustMania/
autorestart=true
startretries=1000

